# elements/neutron-ovs-cleanup/install.d/neutron-source-install/neutron-ovs-cleanup.conf
# neutron-ovs-cleanup
#

start on started openvswitch-switch
stop on runlevel [!2345]

pre-start script
    [ ! -e /etc/neutron/neutron.conf ] && exit 0

    start-stop-daemon --start --chuid neutron --exec /usr/local/bin/neutron-ovs-cleanup -- \
      --config-dir /etc/neutron
end script
