# elements/keepalived/files/keepalived.crontab
*/5 * * * * root [ -s /var/run/keepalived.pid ] && kill -HUP $(cat /var/run/keepalived.pid)
